# rc.d/pfdns
#!/bin/ksh

daemon="/usr/local/bin/pfdns"
daemon_flags="-cfg /etc/pfdns.json"

. /etc/rc.d/rc.subr

rc_start() {
	rc_exec "${daemon} ${daemon_flags} >> /var/log/pfdns.log 2>&1 &"
}

rc_cmd $1

